# hardware/os/meta-goldi-bsp/recipes-bsp/leds/files/set_led_power.service
[Unit]
Description=Sets the "set the Power LED"
After=load-fpga-firmware.service

[Service]
Type=oneshot
ExecStart=/usr/bin/set_led_power

[Install]
WantedBy=multi-user.target