# hardware/os/meta-goldi-distro/recipes-core/boot-and-update/files/update.service
[Unit]
Description=Updates the system

[Service]
Type=simple
Restart=always
RestartSec=60s
ExecStart=/usr/bin/update

[Install]
WantedBy=multi-user.target