# hardware/os/meta-goldi-distro/recipes-core/application/files/goldi-crosslab.service
[Unit]
Description=GOLDi2 Crosslab Client
After=network.target

[Service]
Restart=on-failure
ExecStart=crosslab_client

[Install]
WantedBy=multi-user.target